# docker/dev/supervisord.d/api_gateaway.ini
[program:api_gateway]
command = 
autorestart = true
stderr_logfile = /tmp/api_gateway.log
stdout_logfile = /tmp/api_gateway_std.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10