# docker/dev/supervisord.d/id.ini
[program:id]
command = /go/src/mayo/bin/id.bin
autorestart = true
stderr_logfile = /tmp/id.log
stdout_logfile = /tmp/id_std.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10