# docker/dev/supervisord.d/user.ini
[program:user]
command = /go/src/mayo/bin/api_gateaway.bin
autorestart = true
stderr_logfile = /tmp/api_gateway.log
stdout_logfile = /tmp/api_gateway_std.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10