# docker/dev/supervisord.d/user.ini
[program:user]
command = /go/src/mayo/bin/user.bin
autorestart = true
stderr_logfile = /tmp/user.log
stdout_logfile = /tmp/user_std.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10